# ssh-inhibit-sleep.service
[Unit]
Description=Check for active SSH connections and inhibit sleep
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ssh-inhibit-sleep.sh
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
